# wangge.service.example
[Unit]
Description=网格交易系统
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/wangge
Environment=PATH=/root/wangge/venv/bin
ExecStart=/root/wangge/venv/bin/python /root/wangge/run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
ReadWritePaths=/root/wangge/logs

[Install]
WantedBy=multi-user.target
